# scripts/deploy/stockio.service
[Unit]
Description=Stockio Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=stockio
WorkingDirectory=/opt/stockio
ExecStart=/opt/stockio/.venv/bin/stockio
Restart=always
RestartSec=30

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/stockio/data /opt/stockio/models /opt/stockio/logs
PrivateTmp=true

# Environment
EnvironmentFile=/opt/stockio/.env

[Install]
WantedBy=multi-user.target
