# scripts/deploy/stockio.service
[Unit]
Description=Stockio Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=stockio
Group=stockio
WorkingDirectory=/opt/stockio
ExecStart=/opt/stockio/.venv/bin/stockio run
Restart=always
RestartSec=30

# Security hardening (LXC-compatible — no namespace features)
NoNewPrivileges=true

# Environment
EnvironmentFile=/opt/stockio/.env

[Install]
WantedBy=multi-user.target
