# deploy/stockio-web.service
[Unit]
Description=Stockio Web Dashboard
After=network.target

[Service]
Type=exec
User=stockio
Group=stockio
WorkingDirectory=/opt/stockio
Environment=PATH=/opt/stockio/.venv/bin:/usr/bin:/bin
Environment=HF_HOME=/opt/stockio/data/.hf_cache
EnvironmentFile=/opt/stockio/.env
ExecStart=/opt/stockio/.venv/bin/gunicorn \
    --bind 0.0.0.0:80 \
    --workers 1 \
    --preload \
    --timeout 300 \
    --worker-tmp-dir /tmp \
    --access-logfile /var/log/stockio/access.log \
    --error-logfile /var/log/stockio/error.log \
    "stockio.webapp:app"
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true
ReadWritePaths=/opt/stockio/data /var/log/stockio
ProtectHome=true

[Install]
WantedBy=multi-user.target
